# assist_microphone/rootfs/etc/s6-overlay/s6-rc.d/assist_microphone/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start service
# ==============================================================================
extra_args=()
if bashio::config.true 'debug_logging'; then
    extra_args+=('--debug')
fi

if bashio::config.true 'sound_enabled'; then
    extra_args+=('--snd-command' 'aplay -r 16000 -c 1 -f S16_LE -t raw')
fi

# --- Custom Volume Control via MQTT Discovery ---
if bashio::config.true 'enable_volume_control'; then
    VOLUME_TOPIC=$(bashio::config 'volume_topic') || VOLUME_TOPIC="anker/volume/set"
    STATE_TOPIC="${VOLUME_TOPIC}/state"
    DISCOVERY_TOPIC="homeassistant/number/assist_microphone_custom_volume/config"

    # MQTT Zugangsdaten: erst Supervisor API, dann Add-on Options als Fallback
    MQTT_HOST=""
    MQTT_PORT="1883"
    MQTT_USER=""
    MQTT_PASS=""

    if MQTT_HOST=$(bashio::services mqtt "host" 2>/dev/null) && [[ -n "${MQTT_HOST}" ]]; then
        MQTT_PORT=$(bashio::services mqtt "port" 2>/dev/null) || MQTT_PORT="1883"
        MQTT_USER=$(bashio::services mqtt "username" 2>/dev/null) || MQTT_USER=""
        MQTT_PASS=$(bashio::services mqtt "password" 2>/dev/null) || MQTT_PASS=""
        bashio::log.info "🔊 MQTT via Supervisor API (${MQTT_HOST}:${MQTT_PORT})"
    else
        MQTT_HOST=$(bashio::config 'mqtt_host') || MQTT_HOST="core-mosquitto"
        MQTT_PORT=$(bashio::config 'mqtt_port') || MQTT_PORT="1883"
        MQTT_USER=$(bashio::config 'mqtt_user') || MQTT_USER=""
        MQTT_PASS=$(bashio::config 'mqtt_pass') || MQTT_PASS=""
        bashio::log.info "🔊 MQTT via Add-on Config (${MQTT_HOST}:${MQTT_PORT})"
    fi

    if [[ -z "${MQTT_HOST}" ]]; then
        bashio::log.error "❌ Kein MQTT Host! Bitte mqtt_host in den Add-on Optionen setzen."
    else
        MQTT_ARGS=(-h "${MQTT_HOST}" -u "${MQTT_USER}" -P "${MQTT_PASS}")
        [[ -n "${MQTT_PORT}" ]] && MQTT_ARGS+=(-p "${MQTT_PORT}")

        # MQTT Discovery: HA erstellt automatisch eine number Entity
        DISCOVERY_PAYLOAD=$(jq -n \
            --arg name "Anker Lautstärke" \
            --arg uid "assist_microphone_custom_volume" \
            --arg cmd_t "${VOLUME_TOPIC}" \
            --arg stat_t "${STATE_TOPIC}" \
            '{
                name: $name,
                unique_id: $uid,
                command_topic: $cmd_t,
                state_topic: $stat_t,
                min: 0,
                max: 100,
                step: 5,
                unit_of_measurement: "%",
                icon: "mdi:volume-high",
                mode: "slider"
            }')

        if mosquitto_pub "${MQTT_ARGS[@]}" -r -t "${DISCOVERY_TOPIC}" -m "${DISCOVERY_PAYLOAD}" 2>/dev/null; then
            bashio::log.info "🔊 MQTT Discovery Entity registriert"
        else
            bashio::log.error "❌ MQTT Discovery fehlgeschlagen"
        fi

        # ALSA Mixer erkennen
        ALSA_CONTROL=""
        if amixer sget "Master" &>/dev/null; then
            ALSA_CONTROL="Master"
        else
            ALSA_CONTROL=$(amixer scontrols 2>/dev/null | head -1 | sed "s/.*'\(.*\)'.*/\1/") || true
        fi

        if [[ -z "${ALSA_CONTROL}" ]]; then
            bashio::log.error "❌ Kein ALSA Mixer gefunden! (aplay -l: $(aplay -l 2>&1 || echo 'keine Geräte'))"
        else
            # Initialen State publishen
            CURRENT_VOL=$(amixer sget "${ALSA_CONTROL}" 2>/dev/null | grep -oP '\d+%' | head -1 | tr -d '%') || true
            CURRENT_VOL=${CURRENT_VOL:-50}
            mosquitto_pub "${MQTT_ARGS[@]}" -r -t "${STATE_TOPIC}" -m "${CURRENT_VOL}" 2>/dev/null || true

            bashio::log.info "🔊 Volume Control aktiv: ${VOLUME_TOPIC} → ${ALSA_CONTROL} (aktuell: ${CURRENT_VOL}%)"

            # Auf Lautstärke-Befehle lauschen
            mosquitto_sub "${MQTT_ARGS[@]}" -t "${VOLUME_TOPIC}" | while read -r payload; do
                if [[ "${payload}" =~ ^[0-9]+(\.[0-9]+)?$ ]]; then
                    vol=$(printf '%.0f' "${payload}")
                    if (( vol >= 0 && vol <= 100 )); then
                        bashio::log.info "🔊 Lautstärke → ${vol}%"
                        amixer sset "${ALSA_CONTROL}" "${vol}%" >/dev/null 2>&1 || bashio::log.error "❌ amixer fehlgeschlagen"
                        mosquitto_pub "${MQTT_ARGS[@]}" -r -t "${STATE_TOPIC}" -m "${vol}" 2>/dev/null || true
                    fi
                fi
            done &
        fi
    fi
fi

exec python3 -m wyoming_satellite \
    --name 'assist microphone' \
    --uri 'tcp://0.0.0.0:10700' \
    --mic-command 'arecord -r 16000 -c 1 -f S16_LE -t raw' \
    --snd-command-rate 16000 \
    --awake-wav "$(bashio::config 'awake_wav')" \
    --done-wav "$(bashio::config 'done_wav')" \
    --timer-finished-wav "$(bashio::config 'timer_finished_wav')" \
    --timer-finished-wav-repeat "$(bashio::config 'timer_repeat_count')" "$(bashio::config 'timer_repeat_delay')" \
    --mic-volume-multiplier "$(bashio::config 'mic_volume_multiplier')" \
    --snd-volume-multiplier "$(bashio::config 'sound_volume_multiplier')" \
    --mic-auto-gain "$(bashio::config 'auto_gain')" \
    --mic-noise-suppression "$(bashio::config 'noise_suppression')" \
    --no-zeroconf "${extra_args[@]}"
