# assist_microphone/rootfs/etc/s6-overlay/s6-rc.d/assist_microphone/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start service
# ==============================================================================
extra_args=()
if bashio::config.true 'debug_logging'; then
    extra_args+=('--debug')
fi

if bashio::config.true 'sound_enabled'; then
    extra_args+=('--snd-command' 'aplay -r 16000 -c 1 -f S16_LE -t raw')
fi

# --- START: Custom Volume Control via MQTT Discovery (v1.0.2) ---
if bashio::config.true 'enable_volume_control'; then
    bashio::log.info "🔊 [v1.0.5] Volume Control Block gestartet"

    VOLUME_TOPIC=$(bashio::config 'volume_topic') || VOLUME_TOPIC="anker/volume/set"
    STATE_TOPIC="${VOLUME_TOPIC}/state"
    DISCOVERY_TOPIC="homeassistant/number/assist_microphone_custom_volume/config"
    bashio::log.info "🔊 Topics: cmd=${VOLUME_TOPIC} state=${STATE_TOPIC}"

    # MQTT Zugangsdaten mit Fallback
    MQTT_HOST=$(bashio::services mqtt "host") || MQTT_HOST=""
    MQTT_PORT=$(bashio::services mqtt "port") || MQTT_PORT="1883"
    MQTT_USER=$(bashio::services mqtt "username") || MQTT_USER=""
    MQTT_PASS=$(bashio::services mqtt "password") || MQTT_PASS=""
    bashio::log.info "🔊 MQTT: host=${MQTT_HOST} port=${MQTT_PORT} user=${MQTT_USER}"

    if [[ -z "${MQTT_HOST}" ]]; then
        bashio::log.error "❌ MQTT Host leer! MQTT Service nicht verfügbar?"
    else
        # MQTT Argumente zusammenbauen
        MQTT_ARGS=(-h "${MQTT_HOST}" -u "${MQTT_USER}" -P "${MQTT_PASS}")
        if [[ -n "${MQTT_PORT}" ]]; then
            MQTT_ARGS+=(-p "${MQTT_PORT}")
        fi

        # 1) MQTT Discovery Config veröffentlichen
        bashio::log.info "🔊 Sende Discovery an: ${DISCOVERY_TOPIC}"
        DISCOVERY_PAYLOAD=$(jq -n \
            --arg name "Anker Lautstärke" \
            --arg uid "assist_microphone_custom_volume" \
            --arg cmd_t "${VOLUME_TOPIC}" \
            --arg stat_t "${STATE_TOPIC}" \
            '{
                name: $name,
                unique_id: $uid,
                command_topic: $cmd_t,
                state_topic: $stat_t,
                min: 0,
                max: 100,
                step: 5,
                unit_of_measurement: "%",
                icon: "mdi:volume-high",
                mode: "slider",
                device: {
                    identifiers: ["assist_microphone_custom"],
                    name: "Assist Microphone Custom",
                    manufacturer: "Custom Add-on",
                    model: "USB Audio Volume Control"
                }
            }')
        bashio::log.info "🔊 Discovery Payload: ${DISCOVERY_PAYLOAD}"

        if mosquitto_pub "${MQTT_ARGS[@]}" -r -t "${DISCOVERY_TOPIC}" -m "${DISCOVERY_PAYLOAD}" 2>&1; then
            bashio::log.info "✅ MQTT Discovery veröffentlicht"
        else
            bashio::log.error "❌ mosquitto_pub fehlgeschlagen (exit=$?)"
        fi

        # 2) ALSA Mixer-Control erkennen
        bashio::log.info "🔊 Suche ALSA Mixer..."
        ALSA_CONTROL=""
        if amixer sget "Master" &>/dev/null; then
            ALSA_CONTROL="Master"
        else
            ALSA_CONTROL=$(amixer scontrols 2>/dev/null | head -1 | sed "s/.*'\(.*\)'.*/\1/") || true
        fi
        bashio::log.info "🔊 ALSA Ergebnis: '${ALSA_CONTROL}'"

        if [[ -z "${ALSA_CONTROL}" ]]; then
            bashio::log.error "❌ Kein ALSA Mixer gefunden!"
            bashio::log.error "   aplay -l: $(aplay -l 2>&1 || echo 'fehler')"
            bashio::log.error "   amixer: $(amixer 2>&1 || echo 'fehler')"
        else
            # 3) Initialen State publishen
            CURRENT_VOL=$(amixer sget "${ALSA_CONTROL}" 2>/dev/null | grep -oP '\d+%' | head -1 | tr -d '%') || true
            CURRENT_VOL=${CURRENT_VOL:-50}
            mosquitto_pub "${MQTT_ARGS[@]}" -r -t "${STATE_TOPIC}" -m "${CURRENT_VOL}" 2>&1 || true
            bashio::log.info "🔊 Initialer Wert: ${CURRENT_VOL}%"

            # 4) Auf Lautstärke-Befehle lauschen
            bashio::log.info "🔊 Starte Listener auf: ${VOLUME_TOPIC} (Mixer: ${ALSA_CONTROL})"
            mosquitto_sub "${MQTT_ARGS[@]}" -t "${VOLUME_TOPIC}" | while read -r payload; do
                if [[ "${payload}" =~ ^[0-9]+(\.[0-9]+)?$ ]]; then
                    vol=$(printf '%.0f' "${payload}")
                    if (( vol >= 0 && vol <= 100 )); then
                        bashio::log.info "🔊 Lautstärke → ${vol}%"
                        amixer sset "${ALSA_CONTROL}" "${vol}%" 2>&1 || bashio::log.error "❌ amixer fehlgeschlagen"
                        mosquitto_pub "${MQTT_ARGS[@]}" -r -t "${STATE_TOPIC}" -m "${vol}" 2>&1 || true
                    else
                        bashio::log.warning "⚠️ Wert außerhalb 0–100: ${payload}"
                    fi
                else
                    bashio::log.warning "⚠️ Ungültiger Payload: ${payload}"
                fi
            done &
        fi
    fi
else
    bashio::log.info "🔊 Volume Control ist deaktiviert"
fi
# --- ENDE: Custom Volume Control ---
exec python3 -m wyoming_satellite \
    --name 'assist microphone' \
    --uri 'tcp://0.0.0.0:10700' \
    --mic-command 'arecord -r 16000 -c 1 -f S16_LE -t raw' \
    --snd-command-rate 16000 \
    --awake-wav "$(bashio::config 'awake_wav')" \
    --done-wav "$(bashio::config 'done_wav')" \
    --timer-finished-wav "$(bashio::config 'timer_finished_wav')" \
    --timer-finished-wav-repeat "$(bashio::config 'timer_repeat_count')" "$(bashio::config 'timer_repeat_delay')" \
    --mic-volume-multiplier "$(bashio::config 'mic_volume_multiplier')" \
    --snd-volume-multiplier "$(bashio::config 'sound_volume_multiplier')" \
    --mic-auto-gain "$(bashio::config 'auto_gain')" \
    --mic-noise-suppression "$(bashio::config 'noise_suppression')" \
    --no-zeroconf "${extra_args[@]}"
